# supervisord/ss.conf
[program:ss]
command=sslocal  -c /data/www/shadowsocks.json
stdout_logfile=/data/logs/ss-access.log
stderr_logfile=/data/logs/ss-error.log
autorestart=true
autostart=true
user=root
